# scripts/deploy.service
[Unit]
Description=Rebocap Doc Webhook Server
After=network.target

[Service]
Type=simple
User=nginx
WorkingDirectory=/opt/rebocap-doc
EnvironmentFile=/opt/rebocap-doc/.env
ExecStart=/usr/bin/node /opt/rebocap-doc/scripts/webhook-server.mjs
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/rebocap-doc

[Install]
WantedBy=multi-user.target
