# scripts/deploy.service
[Unit]
Description=Rebocap Doc Webhook Server
After=network.target

[Service]
Type=simple
User=nginx
WorkingDirectory=/opt/rebocap_docs
EnvironmentFile=/opt/rebocap_docs/.env
ExecStart=/usr/bin/node /opt/rebocap_docs/scripts/webhook-server.mjs
Restart=always
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
